# backup/var/lib/dpkg/info/rpi-eeprom.list
/.
/etc
/etc/default
/etc/default/rpi-eeprom-update
/lib
/lib/firmware
/lib/firmware/raspberrypi
/lib/firmware/raspberrypi/bootloader
/lib/firmware/raspberrypi/bootloader/beta
/lib/firmware/raspberrypi/bootloader/beta/pieeprom-2020-07-16.bin
/lib/firmware/raspberrypi/bootloader/beta/pieeprom-2020-07-31.bin
/lib/firmware/raspberrypi/bootloader/beta/recovery.bin
/lib/firmware/raspberrypi/bootloader/beta/vl805-000138a1.bin
/lib/firmware/raspberrypi/bootloader/critical
/lib/firmware/raspberrypi/bootloader/critical/pieeprom-2020-04-16.bin
/lib/firmware/raspberrypi/bootloader/critical/recovery.bin
/lib/firmware/raspberrypi/bootloader/critical/vl805-000137ad.bin
/lib/firmware/raspberrypi/bootloader/release-notes.md
/lib/firmware/raspberrypi/bootloader/stable
/lib/firmware/raspberrypi/bootloader/stable/pieeprom-2020-04-16.bin
/lib/firmware/raspberrypi/bootloader/stable/pieeprom-2020-06-15.bin
/lib/firmware/raspberrypi/bootloader/stable/pieeprom-2020-07-16.bin
/lib/firmware/raspberrypi/bootloader/stable/pieeprom-2020-07-31.bin
/lib/firmware/raspberrypi/bootloader/stable/recovery.bin
/lib/firmware/raspberrypi/bootloader/stable/vl805-000137ad.bin
/lib/firmware/raspberrypi/bootloader/stable/vl805-000138a1.bin
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rpi-eeprom-update.service
/usr
/usr/bin
/usr/bin/rpi-eeprom-config
/usr/bin/rpi-eeprom-update
/usr/share
/usr/share/doc
/usr/share/doc/rpi-eeprom
/usr/share/doc/rpi-eeprom/LICENCE.bootloader
/usr/share/doc/rpi-eeprom/changelog.Debian.gz
/usr/share/doc/rpi-eeprom/copyright
/usr/share/doc/rpi-eeprom/release-notes.md.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/rpi-eeprom-config.1.gz
/usr/share/man/man1/rpi-eeprom-update.1.gz
/var
/var/lib
/var/lib/raspberrypi
/var/lib/raspberrypi/bootloader
/var/lib/raspberrypi/bootloader/backup
